# agent/reach-agent.service
[Unit]
Description=Reach Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/reach-agent
Restart=always
RestartSec=10

# Run as dedicated non-root user
User=reach-agent
Group=reach-agent

# Prevent the process from gaining new privileges via setuid/setgid
NoNewPrivileges=yes

# Isolated /tmp - child processes get their own private temp dir
PrivateTmp=yes

# Protect kernel tunables and control groups from writes
ProtectKernelTunables=yes
ProtectControlGroups=yes

# Drop all capabilities - the agent only needs network access
CapabilityBoundingSet=

Environment=REACH_COMMAND_TIMEOUT_SECONDS=60
Environment=REACH_MAX_OUTPUT_BYTES=50000

[Install]
WantedBy=multi-user.target
